# config/messenger-worker.ini
[program:messenger-consume]
command=php /Users/gopibabus/Personal/symfony-concepts-3/bin/console messenger:consume async_priority_high async --time-limit=3600
user=gopibabus
numprocs=2
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d